# scripts/quickphoto.conf
[program:quickphoto]
command=php artisan queue:listen --tries=2
directory=/var/www
stdout_logfile=/var/www/storage/logs/supervisor.log
redirect_stderr=true
